# daemons/rclone.timer
[Unit]
Description=Run rclone sync periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
